# nginx-php8/s6/php-fpm8/run
#!/bin/sh

exec /usr/sbin/php-fpm81 -R -F --fpm-config /etc/php81/php-fpm.conf
